# debian/prerm
#!/bin/sh
set -e

deb-systemd-invoke stop airtop-fpsvc.service || :
deb-systemd-helper disable airtop-fpsvc.service || :
deb-systemd-invoke stop gpu-thermald.service || :
deb-systemd-helper disable gpu-thermald.service || :

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0